# config/crontab
@reboot (cd /home/ubuntu/karyt/www; ./server.py >> /var/log/karyt/server.log 2>&1)
